# units/api-routerd.socket
[Unit]
Description=api-routerd ListenStream

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target
